# roles/fan_control/templates/fan_control.service.j2
[Unit]
Description=Fan Control Service
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/bin/fan_control.py --monitor
Restart=always

[Install]
WantedBy=multi-user.target
